# sv/sshd/run
#!/usr/bin/env bash

exec /usr/sbin/sshd -D -e 2>&1
