# sv/openvpn/run
#!/usr/bin/env bash

exec /usr/sbin/openvpn --config /etc/openvpn/vpn.ovpn
